# files/system/usr/lib/systemd/system/flatpak-update.timer
[Unit]
Description=Daily Flatpak Updates

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
